# door-bell.service
[Unit]
Description=Door Bell Service

[Service]
ExecStart="/usr/bin/python /home/innay/bell-commander/main.py
User=innay

